# memory-core/mcaas_clean/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; ── Core infrastructure (priority 10) ─────────────────────────────────────────

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres
    -D /data/postgres
    -c listen_addresses=127.0.0.1
    -c port=5432
    -c max_connections=50
    -c shared_buffers=256MB
    -c log_min_messages=warning
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/mc/postgres.log
stderr_logfile=/var/log/mc/postgres.err

[program:redis]
command=redis-server
    --bind 127.0.0.1
    --port 6379
    --maxmemory 256mb
    --maxmemory-policy allkeys-lru
    --save ""
    --appendonly no
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/mc/redis.log
stderr_logfile=/var/log/mc/redis.err

[program:ollama]
command=ollama serve
environment=OLLAMA_HOST="127.0.0.1:11434",OLLAMA_MODELS="/data/ollama"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/mc/ollama.log
stderr_logfile=/var/log/mc/ollama.err

; ── Workers (priority 30) ──────────────────────────────────────────────────────

[program:engram-worker]
command=python -m workers.engram_worker
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=8
stdout_logfile=/var/log/mc/engram.log
stderr_logfile=/var/log/mc/engram.err

[program:eidos-worker]
command=python -m workers.eidos_worker
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=8
stdout_logfile=/var/log/mc/eidos.log
stderr_logfile=/var/log/mc/eidos.err

[program:reflective-scheduler]
command=python -m workers.reflective_scheduler
directory=/app
autostart=true
autorestart=true
priority=40
startsecs=10
stdout_logfile=/var/log/mc/reflective.log
stderr_logfile=/var/log/mc/reflective.err

; ── API (priority 50) ─────────────────────────────────────────────────────────

[program:api]
command=uvicorn api.main:app
    --host 0.0.0.0
    --port 4200
    --workers 2
    --log-level warning
    --access-log
directory=/app
autostart=true
autorestart=true
priority=50
startsecs=12
stdout_logfile=/var/log/mc/api.log
stderr_logfile=/var/log/mc/api.err

; ── Group for clean shutdown ───────────────────────────────────────────────────

[group:workers]
programs=engram-worker,eidos-worker,reflective-scheduler

[group:infrastructure]
programs=postgres,redis,ollama
